# storage_server/utils/backup-purge.cron
@weekly root /usr/bin/backup_purger.sh -d 7 -m 1 -s 4
